# bot.service
[Unit]
Description=Robin Trading Bot
After=network.target

[Service]
# Replace 'username' with your actual GCP username (usually your email prefix)
# Assuming the code is uploaded to /home/username/delta_bot
WorkingDirectory=/home/%u/delta_bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
User=%u
# Environment variables if needed
# Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
